# dot_config/systemd/user/git-autocommit.timer
[Unit]
Description=Run git autocommit hourly

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
