# dot_config/systemd/user/shutdown-prompt.service
[Unit]
Description=Shutdown Prompt Service

[Service]
Type=oneshot
ExecStart=/home/lyka/.local/bin/maybe-shutdown.sh
RemainAfterExit=yes
