# dot_config/systemd/user/sync.service
[Unit]
Description=Run sync.sh script

[Service]
Type=oneshot
ExecStart=/home/lyka/sync.sh
